# etc/supervisor/conf.d/vncd.conf
[program:vncd]
command=bash -c '/dockerstartup/startup.sh'
#/dockerstartup/startup.sh
#sudo  vncserver -localhost  -depth 24 -geometry 1360x768 :1
autorestart=false
directory=/root