# etc/supervisor/conf.d/vncd.conf
[program:vncd_startup]
command=bash -c '/dockerstartup/startup.sh'
autorestart=false
directory=/root