# etc/supervisor/conf.d/vncd.conf
[program:vncd]
command=bash -c 'sudo vncserver -geometry 1280x800 :0'
autorestart=false
directory=/root