# etc/supervisor/conf.d/websockify.conf
[program:websockify]
command=python3 -m websockify --web /usr/share/novnc/utils/../ 6080 localhost:5901
autorestart=true
